# home/.bkps.git.cron
#!/bin/bash

bkpFolders=(
  $HOME/dotfiles
  $HOME/notes
  $HOME/private
  $HOME/code/arafatm/blinkist
)

pushd $(pwd)

for f in ${bkpFolders[@]}; do
  if [[ -d $f ]]; then
    cd $f && echo $(pwd)
    echo "########################################" \
      && git pull \
      && git add . \
      && git commit . -m "bkp" \
      && git push \
      && echo "$(date +%F_%T)" \
      >> $HOME/tmp/cron.log 2>&1
  fi
done

popd
